# systemd/notion-sync.service
[Unit]
Description=Notion Sync Manager Daemon
After=network.target

[Service]
Type=simple
User=feyijimiohioma
Environment=NODE_ENV=production
Environment=NOTION_SYNC_CONFIG_DIR=/home/feyijimiohioma/.notion-sync-manager
WorkingDirectory=/home/feyijimiohioma/projects/Nsma
ExecStart=/home/feyijimiohioma/.nvm/versions/node/v24.11.1/bin/node cli/index.js --daemon
Restart=always
RestartSec=30
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
